# infra/opsd/uk.mattjewell.assistant-opsd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>uk.mattjewell.assistant-opsd</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/mattisfrommars/.local/share/mise/installs/node/22.22.1/bin/node</string>
    <string>/Users/mattisfrommars/opsd/opsd.mjs</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/mattisfrommars/.local/share/mise/shims:/opt/homebrew/bin:/Applications/OrbStack.app/Contents/MacOS/xbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/mattisfrommars</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/mattisfrommars/opsd/opsd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mattisfrommars/opsd/opsd.log</string>
</dict>
</plist>
